# rc.local
#!/bin/sh -e
#
# rc.local
#
# This script is executed at the end of each multiuser runlevel.
# Make sure that the script will "exit 0" on success or any other
# value on error.
#
# In order to enable or disable this script just change the execution
# bits.
#
# By default this script does nothing.

# Keyboard
/usr/bin/loadkeys /home/chip/.config/pocketchip.kmap
# Ralink wifi chipsets
modprobe rt2800usb
echo "148f 5370" > /sys/bus/usb/drivers/rt2800usb/new_id
exit 0
